# webui/gatecrash-webui.service
[Unit]
Description=Gatecrash Web UI
After=network.target

[Service]
Type=simple
User=gatecrash
Group=gatecrash
ExecStart=/opt/gatecrash/webui/venv/bin/python /opt/gatecrash/webui/app.py
WorkingDirectory=/opt/gatecrash/webui
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
